# smdb.service.template
[Unit]
Description=Server Monitoring Discord Bot
After=multi-user.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 <file_path>

[Install]
WantedBy=multi-user.target